# config/systemd/ai-cam-media.service
[Unit]
Description=AICam Media Service
After=network.target

[Service]
Type=simple
User={{TARGET_USER}}
WorkingDirectory={{DEPLOY_PATH}}
ExecStart={{DEPLOY_PATH}}/apps/media_service/target/release/ai-cam-media-service
# The binary catches SIGTERM and tears the pipeline down to NULL
# (releasing the vdevice). This 2 s wait covers the case where the
# handler is missed or the binary crashes mid-shutdown, giving the
# kernel time to release the hailo handle before the next process
# starts and hits HAILO_OUT_OF_PHYSICAL_DEVICES.
ExecStopPost=/bin/sleep 2
Restart=on-failure
RestartSec=5

# Capability grant that lets a producer thread enter realtime
# scheduling without a separate unit change. `CAP_SYS_NICE` grants
# the right to enter realtime; `LimitRTPRIO=20` caps how high it can
# go (well under the kernel softirq/ksoftirqd ceiling at prio 50).
AmbientCapabilities=CAP_SYS_NICE
LimitRTPRIO=20

[Install]
WantedBy=multi-user.target
